# deploy/systemd/recommend-restaurant-tunnel.service
[Unit]
Description=Recommend Restaurant Cloudflare Quick Tunnel
After=network-online.target recommend-restaurant.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared tunnel --url http://127.0.0.1:8011 --no-autoupdate
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
